# lxc/ha-entity-vault.service
[Unit]
Description=HA Entity Vault (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=haev
Group=haev
WorkingDirectory=/srv/ha-entity-vault
EnvironmentFile=-/etc/default/ha-entity-vault
ExecStart=/opt/ha-entity-vault/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/data

[Install]
WantedBy=multi-user.target
